# scripts/tools/pi/scripts/supervisor/ch-storage-client.conf

[program:ch-storage-client]
environment=LD_LIBRARY_PATH="/usr/local/lib"
command=/usr/local/bin/ch-storage-client
